# QRADAR_RULES/QRADAR_RULESSS/05_service_persistence.json
{
  "id": "QRADAR-RULE-005",
  "title": "Suspicious Service / Scheduled-Task Persistence",
  "status": "production",
  "description": "Detects suspicious Windows service (4697) or scheduled task (4698) creation pointing to LOLBINs or temp/appdata paths.",
  "author": "SOC Team",
  "date": "2026/06/24",
  "tags": [
    "attack.persistence",
    "attack.t1543.003",
    "attack.t1053.005"
  ],
  "mitre": {
    "tactic": "Persistence",
    "technique": "T1543",
    "subtechnique": "T1543.003"
  },
  "logsource": {
    "product": "windows",
    "service": "security"
  },
  "qradar": {
    "rule_name": "SOC-005: Suspicious Service / Scheduled-Task Persistence",
    "severity": "HIGH",
    "credibility": 8,
    "relevance": 9,
    "response_action": "create_offense"
  },
  "aql": "SELECT LOGSOURCENAME(logsourceid) AS host, username AS performed_by, DATEFORMAT(starttime, 'yyyy-MM-dd HH:mm:ss') AS event_time, UTF8(payload) AS raw_event FROM events WHERE UTF8(payload) ILIKE '%AgentLogFile=Security%' AND ( UTF8(payload) ILIKE '%EventID=4697%' OR UTF8(payload) ILIKE '%EventID=4698%' ) AND ( UTF8(payload) ILIKE '%powershell%' OR UTF8(payload) ILIKE '%cmd.exe%' OR UTF8(payload) ILIKE '%mshta%' OR UTF8(payload) ILIKE '%rundll32%' OR UTF8(payload) ILIKE '%regsvr32%' OR UTF8(payload) ILIKE '%certutil%' OR UTF8(payload) ILIKE '%\\\\temp\\\\%' OR UTF8(payload) ILIKE '%\\\\appdata\\\\%' OR UTF8(payload) ILIKE '%\\\\users\\\\public%' OR UTF8(payload) ILIKE '%-enc%' OR UTF8(payload) ILIKE '%frombase64%' OR UTF8(payload) ILIKE '%.bat%' OR UTF8(payload) ILIKE '%.vbs%' OR UTF8(payload) ILIKE '%.ps1%' ) ORDER BY starttime DESC LAST 15 MINUTES",
  "falsepositives": [
    "Legitimate software installs",
    "Approved IT automation"
  ],
  "level": "high"
}