# QRADAR_RULES/QRADAR_RULESSS/11_persistence_mechanism_established.json
{
  "id": "QRADAR-RULE-011",
  "title": "Persistence Mechanism Established (Linux)",
  "status": "production",
  "description": "Detects Linux persistence via cron, systemd units, rc.local or authorized_keys modification.",
  "author": "SOC Team",
  "date": "2026/06/24",
  "tags": [
    "attack.persistence",
    "attack.t1053",
    "attack.t1543",
    "attack.t1098.004"
  ],
  "mitre": {
    "tactic": "Persistence",
    "technique": "T1053",
    "subtechnique": "T1098.004"
  },
  "logsource": {
    "product": "linux",
    "service": "auditd"
  },
  "qradar": {
    "rule_name": "SOC-011: Persistence Mechanism Established (Linux)",
    "severity": "HIGH",
    "credibility": 8,
    "relevance": 9,
    "response_action": "create_offense"
  },
  "aql": "SELECT QIDNAME(qid) AS event_name, sourceip AS source_ip, username AS target_user, LOGSOURCENAME(logsourceid) AS log_source, DATEFORMAT(starttime, 'yyyy-MM-dd HH:mm:ss') AS event_time, UTF8(payload) AS raw_event FROM events WHERE LOGSOURCENAME(logsourceid) ILIKE '%LinuxServer%' AND ( UTF8(payload) ILIKE '%/etc/cron%' OR UTF8(payload) ILIKE '%authorized_keys%' OR UTF8(payload) ILIKE '%/etc/systemd/system%' OR UTF8(payload) ILIKE '%/etc/rc.local%' ) ORDER BY starttime DESC LAST 5 MINUTES",
  "falsepositives": [
    "Authorized administrators",
    "Configuration management tools"
  ],
  "level": "high"
}